# docker/supervisord.conf
[supervisord]
nodaemon=true

[program:monitor]
command=/opt/subtle-delay/bin/run-background.sh

